# amplifi-dev/backend/supervisord-celery-ingestion.conf
[supervisord]
nodaemon=true

[program:celery-ingestion]
command=python /app/dist/run_celery.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CELERY_QUEUES="ingestion_queue,document_ingestion_queue",CELERY_POOL="threads",CELERY_CONCURRENCY="2"
